# board/bot/rpi0w/overlay/usr/lib/systemd/system/setup_wifi.service
[Unit]
Before=wpa_supplicant@wlan0.service

[Service]
Type=simple
ExecStart=/usr/local/bin/setup_wpa.sh

[Install]
WantedBy=network.target